# autostart/conky.desktop
[Desktop Entry]
Name=Conky
Comment=Conky
Exec=conky
Terminal=false
Type=Application
Categories=
OnlyShowIn=Awesome;
